# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:filebeat]
command=filebeat -e --strict.perms=false

[program:myapi]
command=dotnet /app/MyApi.dll